# firmware/zbf-firmware.service
[Unit]
Description=zbf-firmware

[Service]
ExecStart=/usr/bin/zbf-firmware

[Install]
WantedBy=multi-user.target
