# firmware/zooblinkofon.service
[Unit]
Description=Das Zooblinkfon

[Service]
ExecStart=/usr/bin/zooblinkofon

[Install]
WantedBy=multi-user.target
